# app-admin/monitor/files/apps-monitor.service
[Unit]
Description=Apps Monitor

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/apps-monitor
